# config/includes.chroot/etc/skel/.config/autostart/rune-ui-setup.desktop
[Desktop Entry]
Type=Application
Exec=/usr/local/bin/rune-ui-setup
Hidden=false
NoDisplay=false
Name=Rune OS UI Setup
X-KDE-autostart-phase=2
