# config/includes.chroot/etc/xdg/autostart/rune-installer.desktop
[Desktop Entry]
Name=Rune Installer
Exec=python3 /opt/rune-installer/installer.py
Type=Application
X-GNOME-Autostart-enabled=true
